# data/init.d/openrc-settingsd.in
#!/sbin/runscript
# Copyright 2012 Alexandre Rostovtsev <tetromino@gentoo.org>
# Distributed under the terms of the GNU General Public License v2

description="localed D-Bus service for OpenRC"

depend() {
	need dbus
}

start() {
	[ -n "${NTP_SERVICE}" ] &&
		OPENRC_SETTINGSD_OPTS="--ntp-service=${NTP_SERVICE} ${OPENRC_SETTINGSD_OPTS}"
	ebegin "Starting openrc-settingsd"
	start-stop-daemon --start --quiet --pidfile "@pidfile@" \
		"@libexecdir@/openrc-settingsd" -- ${OPENRC_SETTINGSD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping openrc-settingsd"
	start-stop-daemon --stop --quiet --pidfile "@pidfile@"
	eend $?
}

# vim: set ft=gentoo-init-d ts=4 :
